# systemd/system/system-service.service
[Unit]
Description=System service for testing.

[Service]
Type=forking
ExecStart=/usr/local/bin/system-service

[Install]
WantedBy=multi-user.target
